# supervisord/list_app/app_web/app_web.ini
[program:dev1_website]
command=npm run start -y
directory=/area/web
autostart=true
autorestart=unexpected
exitcodes=0
priority = 1
user=web