# supervisord/list_app/app_web/nginx.ini
[program:nginx]
command = /usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=unexpected
exitcodes=0
; stdout_logfile=/var/log/nginx/stdout.log
; stderr_logfile=/var/log/nginx/stderr.log